# resources/qol-assist-migration.service.in
[Unit]
Description=QoL Migration Service
ConditionPathExists=@TRACKDIR@/trigger
After=systemd-sysusers.service

[Service]
Type=oneshot
ExecStart=@BINDIR@/qol-assist migrate

[Install]
WantedBy=sysinit.target
